# inittab
# Run god
::sysinit:/sbin/god --init

# Set up the TTY's 1 through 3
tty1:123:askfirst:/sbin/mingetty tty1
tty2:123:respawn:/sbin/mingetty tty2
tty3:123:respawn:/sbin/mingetty tty3

# Stop all services on shutdown
::shutdown:/sbin/god --shutdown

# Killing everything on shutdown
::shutdown:echo :: Sending SIGTERM to all.
::shutdown:/bin/kill -s TERM -1
::shutdown:sleep 1
::shutdown:echo :: Sending SIGKILL to all.
::shutdown:/bin/kill -s KILL -1

# Unmount everything on shutdown
::shutdown:echo :: Unmounting everything.
::shutdown:/bin/umount -a -r
::shutdown:/bin/mount -o remount,ro /
